# infra/prod/vllm-qwen.service
# Ynara · vLLM · Qwen 3.5-9B (agente) — systemd unit de produccion.
#
# Companion de vllm-gemma.service (issue #212). Topologia ADR-012 D2: single_process
# co-residente — gemma 12B + qwen 9B + bge-m3 en la MISMA GPU sin swap. Una unit por
# proceso (ADR-013: N procesos = N units). Levantar el stack completo:
#   systemctl enable --now vllm-gemma.service vllm-qwen.service vllm-bge.service
#
# Flags alineados a ADR-012 / ynara.config.json[llm.serving] / start-vllm.sh:
#   --max-model-len 32768         (ynara.config.json: max_model_len qwen)
#   --gpu-memory-utilization 0.36 (PROVISIONAL, ver disclaimer abajo)
#   --quantization awq_marlin     (ynara.config.json: quantization)
#   --kv-cache-dtype fp8          (ADR-012 D3 / ynara.config.json: kv_cache_dtype)
#   --tool-call-parser hermes     (ADR-009 D2: parser de Qwen)
#   thinking ON configurable (agente, ADR-012 D4): default del modelo, sin flag extra.
#
# DISCLAIMER (#207, mismo que start-vllm.sh): --gpu-memory-utilization y
# --max-model-len son PROVISIONALES, derivados de la medicion en Ollama/GGUF
# (ADR-012). Reconfirmar bajo vLLM/AWQ-Marlin; si los 3 no entran juntos en 16 GB,
# bajar PRIMERO el --max-model-len de Qwen (es el mas grande, 32K).

[Unit]
Description=Ynara vLLM - Qwen 3.5-9B (agente, :8002)
Documentation=https://github.com/Ynara/ynara/blob/main/docs/architecture/adrs/ADR-012-conversational-model-12b-single-process.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ynara
Group=ynara
WorkingDirectory=/opt/ynara

# Paths y cache de HuggingFace en EnvironmentFile externo — NUNCA en la unit
# (regla #4). Ej: HF_HOME, QWEN_MODEL, VLLM_BIN, CUDA_VISIBLE_DEVICES.
EnvironmentFile=/etc/ynara/vllm.env

# --served-model-name = served_name de ynara.config.json. Puerto 8002 = LLM_SERVING
# en .env.example. El agente deja thinking ON (default del modelo, ADR-012 D4): el
# razonamiento paso a paso ayuda a planificar tool calls, no se desactiva.
ExecStart=/usr/bin/env vllm serve ${QWEN_MODEL} \
    --port 8002 \
    --served-model-name qwen \
    --max-model-len 32768 \
    --gpu-memory-utilization 0.36 \
    --quantization awq_marlin \
    --kv-cache-dtype fp8 \
    --enable-auto-tool-choice --tool-call-parser hermes

Restart=on-failure
RestartSec=5
# Carga de pesos ~8,5 s en Ollama/GGUF (ADR-012); margen holgado bajo vLLM/AWQ.
TimeoutStartSec=300
Nice=-5

# Hardening basico. ReadWritePaths: logs + cache de HF.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/ynara /opt/ynara/.cache/huggingface

[Install]
WantedBy=multi-user.target
